# sipXpark/bin/sipxpark.in
#!@BASH@
# sipxpark - Startup script for sipxpark

# chkconfig: 35 85 15
# description: sipxpark is a SIP park for telecommunications
# processname: sipXpark

# Copyright (C) 2007 Pingtel Corp., certain elements licensed under a Contributor Agreement.  
# Contributors retain copyright to elements licensed under a Contributor Agreement.
# Licensed to the User under the LGPL license.

. /etc/rc.d/init.d/functions || exit 1
. @SIPX_LIBEXECDIR@/sipx-utils.sh || exit 1

prog=sipxpark
exec="@SIPX_BINDIR@/${prog}"
config="@SIPX_CONFDIR@/sipxpark-config"
pidfile="@SIPX_RUNDIR@/sipxpark.pid"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

start() {
  [ -x $exec ] || { echo -n "$execError"; echo_failure; echo; exit $ENOEXEC; }
  [ -f $config ] || { echo -n "$invalidFileError"; echo_failure; echo; exit $ENOENT; }
  if [ -n "${NoFork}" ] ; then
    killAnyOrphans @SIPXPBXUSER@ $prog
    runuser - @SIPXPBXUSER@ -c $exec
  else
    checkRunningInstance @SIPXPBXUSER@ $prog > /dev/null
    Status=$?
    [ $Status -eq 0 ] && { echo -n "$prog $runningError"; echo_failure; echo; exit $ERUNNING; } 
    
    echo -n $"Starting $prog: "
    daemon --user @SIPXPBXUSER@ --pidfile $pidfile $exec $pidfile
    
    #check if process was launched in process table
    checkRunningInstanceTimeout @SIPXPBXUSER@ $prog
    Status=$?
    [ $Status -eq 0 ] || { echo_failure; echo; exit $ESRCH; }

    checkIfListening $prog
    Status=$? 
    [ $Status -eq 0 ] || { echo_failure; echo; exit $ENOSCKT; }
    [ $Status -eq 0 ] && { echo_success; echo; }
  fi
}

stop() {
  echo -n $"Stopping sipXpark: "
  killproc "${exec}"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
  killAnyOrphans @SIPXPBXUSER@ $prog
}

restart() {
  stop
  start
}

configtest() {
  logfile="@SIPX_LOGDIR@/sipxpark.log"
  if [ -e $logfile -a ! -w $logfile ]
  then
     echo "Log file '$logfile' exists but is not writable by user '@SIPXPBXUSER@'." >&2
     Status=1
  fi
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    checkRunningInstance @SIPXPBXUSER@ $prog > /dev/null
    Status=$?
    [ $Status -eq 0 ] && restart || :
    ;;
  status)
    status "${prog}"
    Status=$?
    ;;
  configtest)
    configtest
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart|nofork}"
    Status=1
esac

exit $Status
