# sipXpublisher/bin/sipxpublisher.in
#!/bin/bash
# sipxpublisher - Startup script for sipxpublisher

# chkconfig: 35 85 15
# description: sipxpublisher is a SIP MWI server for telecommunications
# processname: sipxpublisher

# Copyright (c) 2012 eZuce, Inc. All rights reserved.
# Contributed to SIPfoundry under a Contributor Agreement

# This software is free software; you can redistribute it and/or modify it under
# the terms of the Affero General Public License (AGPL) as published by the
# Free Software Foundation; either version 3 of the License, or (at your option)
# any later version.

# This software is distributed in the hope that it will be useful, but WITHOUT
# ANY WARRANTY; without even the implied warranty of MERCHANTABILITY or FITNESS
# FOR A PARTICULAR PURPOSE. See the GNU Affero General Public License for more
# details.

. /etc/rc.d/init.d/functions || exit $ENOENT
. @SIPX_LIBEXECDIR@/sipx-utils.sh || { echo_failure; echo; exit $ENOENT; }

prog=sipstatus
progName=sipxpublisher
exec="@SIPX_BINDIR@/${prog}"
config="@SIPX_CONFDIR@/status-config"
pidfile="@SIPX_RUNDIR@/sipxpublisher.pid"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

start() {
  [ -x $exec ] || { echo -n "$execError"; echo_failure; echo; exit $ENOEXEC; }
  [ -f $config ] || { echo -n "$invalidFileError"; echo_failure; echo; exit $ENOENT; }
  if [ -n "${NoFork}" ] ; then
    killAnyOrphans @SIPXPBXUSER@ $exec
    runuser - @SIPXPBXUSER@ -c $exec
  else
    checkRunningInstance @SIPXPBXUSER@ $exec > /dev/null
    Status=$?
    [ $Status -eq 0 ] && { echo -n "$progName $runningError"; echo_failure; echo; exit $ERUNNING; }

    echo -n $"Starting $progName: "
    daemon --user @SIPXPBXUSER@ --pidfile $pidfile $exec $pidfile
    Status=$?

    if [ $# -gt 0 ]; then
      #check if process spawned
      checkRunningInstanceTimeout @SIPXPBXUSER@ $exec
      Status=$?
      [ $Status -eq 0 ] || { echo_failure; echo; exit $ESRCH; }
      checkIfListening $prog
      Status=$?
      [ $Status -eq 0 ] || { echo_failure; echo; exit $ENOSCKT; }
      [ $Status -eq 0 ] && { echo_success; echo; }
    else
      [ $Status -eq 0 ] && { echo_success; echo; }
      [ $Status -eq 0 ] || { echo_failure; echo; }
    fi
  fi
}

stop() {
  echo -n $"Stopping $progName: "
  killproc "${exec}"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
  killAnyOrphans @SIPXPBXUSER@ $exec
}

restart() {
  stop
  start $withCheck
}

asyncrestart() {
  stop
  start
}

configtest() {
  sipx_config_exists @SIPX_CONFDIR@/status-config
  Status=$(($Status+$?))
  sipx_config_exists @SIPX_CONFDIR@/status-plugin.xml
  Status=$(($Status+$?))

  # Check that the log file is writable.
  logfile="@SIPX_LOGDIR@/sipstatus.log"
  if [ -e $logfile -a ! -w $logfile ]
  then
    echo "Log file '$logfile' exists but is not writable by user '@SIPXPBXUSER@'." >&2
    Status=1
  fi
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start $withCheck
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  asyncrestart)
    asyncrestart
    ;;
  condrestart)
    checkRunningInstance @SIPXPBXUSER@ $exec > /dev/null
    Status=$?
    [ $Status -eq 0 ] && restart || :
    ;;
  status)
    checkRunningInstance @SIPXPBXUSER@ $exec $prog
    Status=$?
    ;;
  configtest)
    configtest
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|asyncrestart|reload|force-reload|condrestart|nofork}"
    Status=$EINCORRECTUSAGE
esac

exit $Status
