# sipXpublisher/bin/sipxpublisher.in
#!/bin/bash
# sipxpublisher - Startup script for sipxpublisher

# chkconfig: 35 85 15
# description: sipxpublisher is a SIP MWI server for telecommunications
# processname: sipxpublisher

# Copyright (C) 2012 eZuce Inc., certain elements licensed under a Contributor Agreement.
# Contributors retain copyright to elements licensed under a Contributor Agreement.
# Licensed to the User under the AGPL license.

. /etc/rc.d/init.d/functions || exit 1
. @SIPX_LIBEXECDIR@/sipx-utils.sh || exit 1

prog=sipstatus
exec="@SIPX_BINDIR@/${prog}"
config="@SIPX_CONFDIR@/status-config"
pidfile="@SIPX_RUNDIR@/sipxpublisher.pid"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

start() {
    [ -x $exec ] || exit 5
    [ -f $config ] || exit 6
    echo -n $"Starting $prog: "
    # if not running, start it up here, usually something like "daemon $exec"
    daemon --user @SIPXPBXUSER@ --pidfile $pidfile $exec $pidfile
    Status=$?
}

stop() {
  echo -n $"Stopping sipxpublisher: "
  killproc -p "${pidfile}" -t30 -TERM "${binfile}"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
}

restart() {
  stop
  start
}

configtest() {
  sipx_config_exists @SIPX_CONFDIR@/status-config
  Status=$(($Status+$?))
  sipx_config_exists @SIPX_CONFDIR@/status-plugin.xml
  Status=$(($Status+$?))

  # Check that the log file is writable.
  logfile="@SIPX_LOGDIR@/sipstatus.log"
  if [ -e $logfile -a ! -w $logfile ]
  then
    echo "Log file '$logfile' exists but is not writable by user '@SIPXPBXUSER@'." >&2
    Status=1
  fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    [ -f ${pidfile} ] && restart || :
    ;;
  status)
    status -p ${pidfile} ${exec}
    Status=$?
    ;;
  configtest)
    configtest
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
    Status=1
esac

exit $Status
