# sipXpage/bin/sipxpage.in
#!/bin/bash
# sipxpage - Startup script for sipxpage

# chkconfig: 35 85 15
# description: sipxpage is a SIP proxy for telecommunications
# processname: sipxpage

# Copyright (C) 2007 Pingtel Corp., certain elements licensed under a Contributor Agreement.  
# Contributors retain copyright to elements licensed under a Contributor Agreement.
# Licensed to the User under the LGPL license.

. /etc/rc.d/init.d/functions || exit 1
. @SIPX_LIBEXECDIR@/sipx-utils.sh || exit 1

prog=sipxpage
pidfile="@SIPX_RUNDIR@/sipxpage.pid"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

checkRunningInstance() {
  if [ -f "${pidfile}" ]; then
    pid=`cat ${pidfile}`
    if checkpid $pid 2>&1; then
      echo "Process $prog is already running with pid $pid"
      exit 1
    fi
  fi
}

start() {
    checkRunningInstance
    JavaCmd=`@SIPX_BINDIR@/sipx-config --java`
    Dependencies=`@SIPX_BINDIR@/java-dep -d @SIPX_JAVADIR@/sipXcommons @page_PKGS@`
    export CLASSPATH=`echo @SIPX_JAVADIR@/sipXpage/*.jar ${Dependencies} | sed -e 's/ /:/g'`

    # -Dprocname is there so cfengine can find it
    Command="$JavaCmd \
      -Dprocname=sipxpage \
      -Dconf.dir=@SIPX_CONFDIR@ \
      ${SIPXPAGE_OPTS} \
      org.sipfoundry.sipxpage.SipXpage \
      $Args"
    if [ -n "${NoFork}" ] ; then
      runuser -s /bin/bash @SIPXPBXUSER@ -c "${Command}"
    else
      echo -n $"Starting $prog: "
      runuser -s /bin/bash @SIPXPBXUSER@ -c "${Command}" >/dev/null 2>&1 &
      echo $! > ${pidfile}
      echo_success
    fi
}

stop() {
  echo -n $"Stopping sipxpage: "
  killproc -p "${pidfile}" "sipxpage"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
}

restart() {
  stop
  start
}

configtest() {
  sipx_config_exists @SIPX_CONFDIR@/sipxpage.properties 
  Status=$(($Status+$?))

  # check validity of xml files
  @bindir@/sipx-validate-xml @SIPX_CONFDIR@/autoattendants.xml
  Status=$(($Status+$?))

  # Check that the log file is writable.
  logfile="@SIPX_LOGDIR@/sipxpage.log"
  if [ -e $logfile -a ! -w $logfile ]
  then
      echo "Log file '$logfile' exists but is not writable." >&2
      Status=1
  fi
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    [ -f ${pidfile} ] && restart || :
    ;;
  status)
    status -p ${pidfile} java
    ;;
  configtest)
    configtest
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart|nofork}"
    Status=1
esac

exit $Status
