# sipXsupervisor/bin/sipxsupervisor.in
#!/bin/bash
#
# sipxsupervisor         Start/Stop the cfengine server daemon.
#
# chkconfig: - 96 20
# description: sipxsupervisor is responsible from giving out configuration files to
#              those cfengine clients, who wish to update their configs.
#              For more information, visit http://www.cfengine.org
# processname: cf-serverd
# config: /var/cfengine/masterfiles

# Source function library.
. /etc/init.d/functions

config="@SIPX_CFINPUTS@/sipxsupervisord.cf"
exec="/usr/sbin/cf-serverd"
prog=$(basename $exec)
# used to tell which config is in use.
configid=/var/run/sipxsupervisor-configid

[ -e /etc/sysconfig/sipxsupervisor ] && . /etc/sysconfig/sipxsupervisor

lockfile=/var/lock/subsys/$prog
sipxagentlockfile="@SIPX_RUNDIR@/sipxagent.lock"

start() {
	echo -n $config > $configid
	echo -n $"Starting $prog: "
	daemon $exec -f $config
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch $lockfile

	#if sipxagent is not running make sure the host is pulling the last configuration on startup.
    #this applies in cases when secondary host starts after master, and master changed configuration
    #while the secondary is down. sipxagent has to be run as root because we don't want to go through
    #supervisor again
	if [ ! -f $sipxagentlockfile ]
	then
		runuser - root -s @SIPX_BINDIR@/sipxagent
	fi		

	return $RETVAL
}

nofork() {
	$exec -v -f $config
}

stop() {
	echo -n $"Stopping $prog: "
	killproc $prog
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f $lockfile
	return $RETVAL
}

restart() {
  	stop
	start
}

reload() {
    restart
}

force_reload() {
    restart
}

fdrstatus() {
	status $prog
}

case "$1" in
    start|stop|restart|reload|nofork)
	$1
	;;    
    force-reload)
        force_reload
        ;;
    status)
  	fdrstatus
	;;
    # condrestart - only restarts if service *is* running
    condrestart)
        [ ! -f $lockfile ] || restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|condstart|reload|force-reload|nofork}"
        exit 3
esac
exit $?
