# sipXopenacd/bin/sipxopenacd.in
#!/bin/bash
# openacd - Startup script for openacd

# chkconfig: 35 85 15
# description: openacd is a call center application based on FreeSWITCH
# processname: openacd


. /etc/rc.d/init.d/functions || exit 1
. @SIPX_LIBEXECDIR@/sipx-utils.sh || exit 1

prog=openacd
exec="@OPENACD_BINDIR@/${prog}"

[ -e @SIPX_CONFDIR@/${prog}rc ] && . @SIPX_CONFDIR@/${prog}rc

erlmongo_setup() {
  target=erlmongo

  for d in /usr/local/lib /usr/lib64 /usr/lib; do
    candidate=`readlink -f $d/erlang/lib/$target-*`
    if test -d $candidate/ebin; then
      ln -snf $candidate @OPENACD_DIR@/plugin.d/erlmongo
      return
    fi
  done
  echo "$target library not found" 1>&2
  exit 1
}

setup() {
  erlmongo_setup

  # HACK!. What we want is to run openacd as the user
  # the rpm/cfengine script created.
  chown -R @SIPXPBXUSER@:@SIPXPBXGROUP@ @OPENACD_DIR@
}

start() {
  setup
  [ -x $exec ] || exit 5
  if [ -n "${NoFork}" ] ; then
    runuser - @SIPXPBXUSER@ -c "$exec console"
  else
    echo -n $"Starting $prog: "
    runuser - @SIPXPBXUSER@ -c "export ERL_LIBS=$ERL_LIBS:@OPENACD_DIR@/plugin.d:/@OPENACD_DIR@/lib; $exec start -noinput"
    Status=$?
    echo_success
  fi
}

stop() {
  echo -n $"Stopping openacd: "
  $exec stop
  Status=$?
  if [ $Status -eq 0 ]; then
    echo_success
  else
    echo_failure
  fi
  echo
}

restart() {
  stop
  start
}

resync() {
  dir=@OPENACD_DBDIR@/run/Mnesia.openacd@`hostname -f`
  ! test -d $dir || rm -rf $dir
  Status=$?
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  status)
    runuser - @SIPXPBXUSER@ -c "$exec ping"
    Status=$?
    ;;
  configtest)
    configtest
    ;;
  resync)
    resync
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|resync|nofork}"
    Status=1
esac

exit $Status
