# sipXrls/bin/sipxrls.in
#!@BASH@
# sipxrls - Startup script for sipxrls

# chkconfig: 35 98 15
# description: sipxrls is a SIP rls for telecommunications
# processname: sipXrls

# Copyright (C) 2007 Pingtel Corp., certain elements licensed under a Contributor Agreement.
# Contributors retain copyright to elements licensed under a Contributor Agreement.
# Licensed to the User under the LGPL license.

. /etc/rc.d/init.d/functions || exit $ENOENT
. @SIPX_LIBEXECDIR@/sipx-utils.sh || { echo_failure; echo; exit $ENOENT; }

prog=sipxrls
exec="@SIPX_BINDIR@/${prog}"
config="@SIPX_CONFDIR@/sipxrls-config"
pidfile="@SIPX_RUNDIR@/sipxrls.pid"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

start() {
  [ -x $exec ] || { echo -n "$execError"; echo_failure; echo; exit $ENOEXEC; }
  [ -f $config ] || { echo -n "$invalidFileError"; echo_failure; echo; exit $ENOENT; }
  if [ -n "${NoFork}" ] ; then
    killAnyOrphans @SIPXPBXUSER@ $exec
    runuser - @SIPXPBXUSER@ -c $exec
  else
    checkRunningInstance @SIPXPBXUSER@ $exec $prog > /dev/null
    Status=$?
    [ $Status -eq 0 ] && { echo -n "$prog $runningError"; echo_failure; echo; exit $ERUNNING; }

    echo -n $"Starting $prog: "
    daemon --user @SIPXPBXUSER@ --pidfile $pidfile $exec $pidfile
    Status=$?

    if [ $# -gt 0 ]; then
      #check if process spawned
      checkRunningInstanceTimeout @SIPXPBXUSER@ $exec
      Status=$?
      [ $Status -eq 0 ] || { echo_failure; echo; exit $ESRCH; }
      checkIfListening $prog
      Status=$?
      [ $Status -eq 0 ] || { echo_failure; echo; exit $ENOSCKT; }
      [ $Status -eq 0 ] && { echo_success; echo; }
    else
      [ $Status -eq 0 ] && { echo_success; echo; }
      [ $Status -eq 0 ] || { echo_failure; echo; }
    fi
  fi
}

stop() {
  echo -n $"Stopping $prog: "
  killproc "${exec}"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
  killAnyOrphans @SIPXPBXUSER@ $exec
}

restart() {
  stop
  start $withCheck
}

asyncrestart() {
  stop
  start
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start $withCheck
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  asyncrestart)
    asyncrestart
    ;;
  condrestart)
    checkRunningInstance @SIPXPBXUSER@ $exec > /dev/null
    Status=$?
    [ $Status -eq 0 ] && restart || :
    ;;
  status)
    checkRunningInstance @SIPXPBXUSER@ $exec $prog
    Status=$?
    ;;
  configtest)
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|asyncrestart|reload|force-reload|condrestart|nofork}"
    Status=$EINCORRECTUSAGE
esac

exit $Status
